# packaging/systemd/lianli-daemon.service
[Unit]
Description=Lian Li Device Daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/lianli-daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
